# setup/resources/partial/crontab/user

# ######################################################################
#
# Section added by: clvLabs/dotfiles
#

# Save a new Chuck Norris quote in motd every day
0 0 * * * /home/tony/bin/cnquote > /etc/motd

# ######################################################################
