# setup/resources/crontab

# Add a Chuck Norris quote in motd
0 0 * * * /home/tony/bin/cnquote > /etc/motd
